# templates/onepam-agent.service.epp
<%- | String $install_dir | -%>
[Unit]
Description=OnePAM Agent
Documentation=https://onepam.com/docs
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=notify
User=root
Group=root
EnvironmentFile=-<%= $install_dir %>/etc/agent.env
ExecStart=<%= $install_dir %>/bin/onepam-agent \
    --server=${AGENT_API_URL} \
    --tenant-id=${AGENT_TENANT_ID} \
    --agent-id=${AGENT_ID} \
    --data-dir=<%= $install_dir %>/data \
    --log-level=${AGENT_LOG_LEVEL}
Restart=always
RestartSec=10
WatchdogSec=120
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=onepam-agent

[Install]
WantedBy=multi-user.target
